# systemd/fourmore-osm-update.service
[Unit]
Description=FourMore OSM Data Update
Documentation=https://github.com/yourusername/fourmore
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=YOUR_USERNAME
Group=YOUR_GROUP

# Set the working directory to your FourMore project
WorkingDirectory=/path/to/fourmore

# Load environment variables if you have a .env.local file
EnvironmentFile=-/path/to/fourmore/.env.local

# Run the update via make
ExecStart=/usr/bin/make db-update

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fourmore-osm-update

# Resource limits (optional)
# Prevent runaway updates from consuming all resources
MemoryMax=2G
CPUQuota=80%

# Restart policy
Restart=on-failure
RestartSec=5min

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
