# systemd/fourmore-osm-update.timer
[Unit]
Description=FourMore OSM Data Update Timer
Documentation=https://github.com/mvexel/osm-fourmore/systemd
Requires=fourmore-osm-update.service

[Timer]
# Run every Sunday at 2:00 AM
OnCalendar=Sun *-*-* 02:00:00

# If the system was off during the scheduled time, run on next boot
Persistent=true

# Add randomization to avoid load spikes (0-30 minutes)
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
